# backend/supervisor.conf
[unix_http_server]
file=./supervisor.sock
chmod=0700
chown=0:0

[inet_http_server]
port=0.0.0.0:8888
username=admin
password=123123

[supervisord]
logfile=./logs/supervisor.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=./logs/supervisor.pid
nodaemon=false
minfds=1024
minprocs=200

[supervisorctl]
serverurl=http://0.0.0.0:8888

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:rq_worker_1]
command=./up worker
autostart=true
startsecs=5
autorestart=true
startretries=3
user=root
priority=999
redirect_stderr=true
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=20
stdout_logfile=./logs/rq_worker.log
stopasgroup=false
killasgroup=false

[program:rq_worker_2]
command=./up worker
autostart=true
startsecs=5
autorestart=true
startretries=3
user=root
priority=999
redirect_stderr=true
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=20
stdout_logfile=./logs/rq_worker.log
stopasgroup=false
killasgroup=false

[program:rq_worker_3]
command=./up worker
autostart=true
startsecs=5
autorestart=true
startretries=3
user=root
priority=999
redirect_stderr=true
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=20
stdout_logfile=./logs/rq_worker.log
stopasgroup=false
killasgroup=false

[program:rq_worker_4]
command=./up worker
autostart=true
startsecs=5
autorestart=true
startretries=3
user=root
priority=999
redirect_stderr=true
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=20
stdout_logfile=./logs/rq_worker.log
stopasgroup=false
killasgroup=false

[program:rq_worker_5]
command=./up worker
autostart=true
startsecs=5
autorestart=true
startretries=3
user=root
priority=999
redirect_stderr=true
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=20
stdout_logfile=./logs/rq_worker.log
stopasgroup=false
killasgroup=false

[program:rq_worker_6]
command=./up worker
autostart=true
startsecs=5
autorestart=true
startretries=3
user=root
priority=999
redirect_stderr=true
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=20
stdout_logfile=./logs/rq_worker.log
stopasgroup=false
killasgroup=false

[program:rq_worker_7]
command=./up worker
autostart=true
startsecs=5
autorestart=true
startretries=3
user=root
priority=999
redirect_stderr=true
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=20
stdout_logfile=./logs/rq_worker.log
stopasgroup=false
killasgroup=false

[program:rq_worker_8]
command=./up worker
autostart=true
startsecs=5
autorestart=true
startretries=3
user=root
priority=999
redirect_stderr=true
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=20
stdout_logfile=./logs/rq_worker.log
stopasgroup=false
killasgroup=false

[program:rq_worker_9]
command=./up worker
autostart=true
startsecs=5
autorestart=true
startretries=3
user=root
priority=999
redirect_stderr=true
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=20
stdout_logfile=./logs/rq_worker.log
stopasgroup=false
killasgroup=false

[program:rq_worker_10]
command=./up worker
autostart=true
startsecs=5
autorestart=true
startretries=3
user=root
priority=999
redirect_stderr=true
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=20
stdout_logfile=./logs/rq_worker.log
stopasgroup=false
killasgroup=false

[program:rq_worker_11]
command=./up worker
autostart=true
startsecs=5
autorestart=true
startretries=3
user=root
priority=999
redirect_stderr=true
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=20
stdout_logfile=./logs/rq_worker.log
stopasgroup=false
killasgroup=false

[program:rq_worker_12]
command=./up worker
autostart=true
startsecs=5
autorestart=true
startretries=3
user=root
priority=999
redirect_stderr=true
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=20
stdout_logfile=./logs/rq_worker.log
stopasgroup=false
killasgroup=false

[program:rq_scheduler]
command=./up scheduler
autostart=true
startsecs=5
autorestart=true
startretries=3
user=root
priority=999
redirect_stderr=true
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=20
stdout_logfile=./logs/rq_scheduler.log
stopasgroup=false
killasgroup=false

[program:zk_watcher]
command=.env/bin/python  ./manage.py script --watch-zks
autostart=true
startsecs=5
autorestart=true
startretries=3
user=root
priority=999
redirect_stderr=true
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=20
stdout_logfile=./logs/zk_watcher.log
stopasgroup=false
killasgroup=false

